# supervisord/conf.d/myweb.conf
[program:mywebapi]
command=/example/example_gunicorn  # monitor target
user=quokka
audostart=true
autorestart=true
