# Deployment/supervisord/conf.d/myweb.conf
[program:my-nginx]
command=/usr/sbin/nginx -g 'daemon off;'
user=www-data
audostart=true
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
